# etc/cron.d/pyropossum
* * * * * root /etc/init.d/pyropossum start >> /var/log/pyropossum.cron
